# package-data/cron.d
13 * * * * root /usr/sbin/rassh cleanup > /dev/null